# roles/configure/templates/leihs_create-master-secret.path
[Path]
PathChanged={{LEIHS_ROOT_DIR}}/config/
PathChanged={{LEIHS_ROOT_DIR}}/config/master-secret.txt
PathModified={{LEIHS_ROOT_DIR}}/config/master-secret.txt
PathExists={{LEIHS_ROOT_DIR}}/config/master-secret.txt
